# redis/supervisord/supervisord-redis.conf
[program:redis]
user=redis
command=/usr/bin/redis-server --protected-mode no
numprocs=1
autostart=true
autorestart=true
