# mariadb/supervisord/supervisord-mariadb.conf
[program:mysqld]
command=/usr/sbin/mysqld
numprocs=1
autostart=true
autorestart=true
